# etcs/autostart/background.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Wallpaper Downloader
Comment=Download Wallpaper
Exec=background.sh
