# etcs/autostart/gpaste.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Gpaste
Comment=Set Gpaste Daemon
Exec=gpaste.sh
