# etcs/autostart/polybar.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=polybar
Comment=Launch Polybar
Exec=polybar.sh
