# etcs/autostart/xrdb.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Xrdb
Comment=Set Xrdb
Exec=xrdb.sh
