# etcs/autostart/xset.desktop
[Desktop Entry]
Icon=utilities-terminal
Type=Application
Name=Xrandr
Comment=Set Xrandr
Exec=xset.sh
